# root/etc/init.d/telemt
#!/bin/sh /etc/rc.common
START=99
USE_PROCD=1
PROG=/usr/bin/telemt
CONF=/var/etc/telemt.toml

start_service() {
    config_load telemt
    local enabled mode port domain use_socks socks_addr socks_user socks_pass
    local tm_handshake tm_connect tm_keepalive tm_ack
    local use_mp use_stun ad_tag announce_ip prefer_ipv6
    
    config_get_bool enabled general enabled 0
    [ "$enabled" -eq 1 ] || return 0
    
    config_get mode general mode "tls"
    config_get port general port "4443"
    config_get domain general domain "google.com"
    config_get_bool use_socks general use_socks 0
    config_get socks_addr general socks_addr "127.0.0.1:2080"
    config_get socks_user general socks_user ""
    config_get socks_pass general socks_pass ""
    
    config_get tm_handshake general tm_handshake "15"
    config_get tm_connect general tm_connect "10"
    config_get tm_keepalive general tm_keepalive "60"
    config_get tm_ack general tm_ack "300"
    
    config_get_bool use_mp general use_middle_proxy 0
    config_get_bool use_stun general use_stun 0
    config_get ad_tag general ad_tag ""
    config_get announce_ip general announce_ip ""
    config_get_bool prefer_ipv6 general prefer_ipv6 0
    
    mkdir -p /var/etc
    
    local mp_str="false"; [ "$use_mp" -eq 1 ] && mp_str="true"
    local stun_str="false"; [ "$use_stun" -eq 1 ] && stun_str="true"
    local ipv6_str="false"; [ "$prefer_ipv6" -eq 1 ] && ipv6_str="true"

    local m_cls="false"; local m_sec="false"; local m_tls="false"
    case "$mode" in
        "classic") m_cls="true" ;;
        "dd") m_sec="true" ;;
        "all") m_cls="true"; m_sec="true"; m_tls="true" ;;
        *) m_tls="true" ;;
    esac

    cat <<TOML > "$CONF"
show_link = []

[general]
fast_mode = true
use_middle_proxy = $mp_str
stun_probing = $stun_str
prefer_ipv6 = $ipv6_str

[general.modes]
classic = $m_cls
secure = $m_sec
tls = $m_tls
TOML

    [ -n "$ad_tag" ] && echo "ad_tag = \"$ad_tag\"" >> "$CONF"

    cat <<TOML >> "$CONF"

[server]
port = $port
listen_addr_ipv4 = "0.0.0.0"
[[server.listeners]]
ip = "0.0.0.0"
TOML

    [ -n "$announce_ip" ] && echo "announce_ip = \"$announce_ip\"" >> "$CONF"
    
    if [ "$prefer_ipv6" -eq 1 ]; then
        echo "[[server.listeners]]" >> "$CONF"
        echo "ip = \"::\"" >> "$CONF"
    fi

    cat <<TOML >> "$CONF"

[timeouts]
client_handshake = $tm_handshake
tg_connect = $tm_connect
client_keepalive = $tm_keepalive
client_ack = $tm_ack

[censorship]
tls_domain = "$domain"
mask = true
mask_port = 443

[access]
replay_check_len = 65536
TOML

    rm -f /tmp/telemt_users.tmp /tmp/telemt_conns.tmp /tmp/telemt_quota.tmp

    config_foreach add_user_secret user
    config_foreach add_user_conns user
    config_foreach add_user_quota user

    echo -e "\n[access.users]" >> "$CONF"
    if [ -s /tmp/telemt_users.tmp ]; then cat /tmp/telemt_users.tmp >> "$CONF"; fi

    if [ -s /tmp/telemt_conns.tmp ]; then
        echo -e "\n[access.user_max_tcp_conns]" >> "$CONF"
        cat /tmp/telemt_conns.tmp >> "$CONF"
    fi

    if [ -s /tmp/telemt_quota.tmp ]; then
        echo -e "\n[access.user_data_quota]" >> "$CONF"
        cat /tmp/telemt_quota.tmp >> "$CONF"
    fi

    rm -f /tmp/telemt_users.tmp /tmp/telemt_conns.tmp /tmp/telemt_quota.tmp
    
    cat <<TOML >> "$CONF"

[dc_overrides]
"203" = "91.105.192.100:443"
TOML

    local up_str="[[upstreams]]\ntype = \"direct\"\nenabled = true\nweight = 10"
    if [ "$use_socks" -eq 1 ]; then
        up_str="[[upstreams]]\ntype = \"socks5\"\naddress = \"$socks_addr\"\nenabled = true\nweight = 1"
        [ -n "$socks_user" ] && [ -n "$socks_pass" ] && up_str="$up_str\nusername = \"$socks_user\"\npassword = \"$socks_pass\""
    fi
    echo -e "$up_str" >> "$CONF"
    
    procd_open_instance
    procd_set_param command "$PROG" "$CONF"
    procd_set_param file "$CONF"
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_set_param respawn
    procd_close_instance
}

add_user_secret() {
    local secret; config_get secret "$1" secret
    [ -n "$secret" ] && echo "$1 = \"$secret\"" >> /tmp/telemt_users.tmp
}
add_user_conns() {
    local conns; config_get conns "$1" max_tcp_conns
    [ -n "$conns" ] && [ "$conns" != "unlimited" ] && echo "$1 = $conns" >> /tmp/telemt_conns.tmp
}
add_user_quota() {
    local quota; config_get quota "$1" data_quota
    if [ -n "$quota" ] && [ "$quota" != "unlimited" ]; then
        local bytes=$(echo "$quota" | sed 's/,/\./g' | awk '{printf "%d", $1 * 1073741824}')
        [ -n "$bytes" ] && [ "$bytes" -gt 0 ] && echo "$1 = $bytes" >> /tmp/telemt_quota.tmp
    fi
}

service_triggers() { procd_add_reload_trigger "telemt"; }
